# overlay/opt/lebot/lib/cron.d/lebot
28 23 5,19 * * ID=letsencrypt-update /update-certificates.sh
